# digpig-hello.service
[Unit]
Description=DigPig Hello on Center button (Waveshare 1.44" ST7735S)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/victor/lcd-test
# Run with the project's virtualenv so Python packages are found
Environment="VIRTUAL_ENV=/home/victor/lcd-test/venv"
Environment="PATH=/home/victor/lcd-test/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/home/victor/lcd-test/venv/bin/python /home/victor/lcd-test/hello-on-center.py
Restart=on-failure
User=victor
Group=victor

# Give access to GPIO/SPI
AmbientCapabilities=CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
