# digpig-hello.service
[Unit]
Description=DigPig Hello (compat wrapper) — combined into web server
After=network-online.target
Wants=network-online.target digpig-web.service
PartOf=digpig-web.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Start/ensure the combined web service is running; this service does nothing itself
ExecStart=/bin/true
User=victor
Group=victor

[Install]
WantedBy=multi-user.target
